# pacman-db-sync.service
[Unit]
Description=Sync pacman database

[Service]
ExecStart=/opt/pacman-db-sync/pacman-db-sync.sh
